# code/modules/boxroles/files/znc.service
[Unit]
Description=ZNC IRC bouncer

[Service]
Type=simple
ExecStart=/usr/bin/znc -f
User=bouncer

[Install]
WantedBy=default.target
